# cron.monthly/musage
#!/bin/sh

btrfs balance start -musage=5 /
