# cron.weekly/dusage
#!/bin/sh

btrfs balance start -dusage=10 /
